# deploy/systemd/openclaw-sidecar.service
[Unit]
Description=OpenClaw 3-Agent Sidecar
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/openclaw-3agent-sidecar
EnvironmentFile=/opt/openclaw-3agent-sidecar/.env
ExecStart=/opt/openclaw-3agent-sidecar/.venv/bin/openclaw-sidecar
Restart=always
RestartSec=3

# Optional hardening knobs
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
